# crontab/packagistRefresh.sh
#!/bin/bash
###############################################################
# FileName: packagistRefresh.sh
# Author: Chao Hong <teddy.hongchao@gmail.com>
# Program:
# History:Sat 03 Mar 2018 05:40:55 PM CST
###############################################################

# php /home/teddy/packagist/statistics.php >> /home/teddy/packagist/v.log
# php /home/teddy/packagist/statistics.php >> /dev/null
