# crontab/searchNewPackages.sh
#!/bin/bash
###############################################################
# FileName: searchNewPackages.sh
# Author: Chao Hong <teddy.hongchao@gmail.com>
# Program:
# History:Sat 03 Mar 2018 05:40:55 PM CST
###############################################################

php $HOME/github/packagist-rank/src/PackageList.php >> /dev/null
